# crontab
*/3 * * * * /cs/usr/cppguy/safe/goodnews/schedule_run.sh